# app-admin/chef-expander/files/initd/chef-expander
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net chef-solr
}

CHEF_EXPANDER_BINARY=${CHEF_EXPANDER_BINARY:-/usr/bin/chef-expander}
CHEF_EXPANDER_PIDFILE=${CHEF_EXPANDER_PIDFILE:-/var/run/chef/expander.pid}
CHEF_EXPANDER_LOGFILE=${CHEF_EXPANDER_LOGFILE:-/var/log/chef/expander.log}
CHEF_EXPANDER_NODE_COUNT=${CHEF_EXPANDER_NODE_COUNT:-3}
CHEF_USER=${CHEF_USER:-chef}
CHEF_GROUP=${CHEF_GROUP:-chef}

CHEF_DIRS=( /etc/chef /var/lib/chef /var/log/chef /var/run/chef )

fix_owners() {
	for i in ${CHEF_DIRS[@]} ${CONFIGFILE}; do
		chown ${CHEF_USER}:${CHEF_GROUP} ${i}
	done
}

start() {
	ebegin "Starting Chef expander"
	fix_owners
	start-stop-daemon --start --quiet --chuid ${CHEF_USER}:${CHEF_GROUP} \
		--pidfile ${CHEF_EXPANDER_PIDFILE} --exec ${CHEF_EXPANDER_BINARY} \
		-- -d -n ${CHEF_EXPANDER_NODE_COUNT} -P ${CHEF_EXPANDER_PIDFILE} -L ${CHEF_EXPANDER_LOGFILE} \
		${CHEF_EXPANDER_OPTS}
	eend $?
}

stop() {
	ebegin "Shutting down Chef expander"
	start-stop-daemon --stop --quiet --pidfile ${CHEF_EXPANDER_PIDFILE}
	eend $?
}
